# supervisord.conf
; Sample supervisor config file.

[unix_http_server]
file=/run/supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; (ip_address:port specifier, *:port for all iface)

[supervisord]
logfile=/var/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true               ; (start in foreground if true;default false)
childlogdir=/var/log/supervisor ; ('AUTO' child log dir, default $TEMP)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:mysqld]
command=/usr/bin/mysqld_safe --datadir='/var/lib/mysql/' --user=root
directory=/usr
autostart=true
autorestart=unexpected
startsecs=10
startretries=0

[program:nginx]
command=/usr/bin/nginx -c /etc/nginx/nginx.conf
autostart=true
autorestart=unexpected
startsecs=1
startretries=0

[program:phpfpm]
command=/usr/bin/php-fpm -R -F -c /etc/php/ -y /etc/php/php-fpm.conf
autostart=true
autorestart=unexpected
startsecs=1
startretries=0
